# jenkins.conf
[program:jenkins]
command=java -jar /usr/share/jenkins/jenkins.war --prefix=$JENKINS_PREFIX
redirect_stderr=true